# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx

[program:node]
command=node server/chat.js

